# extensions/elk/ansible/roles/rita/files/Rita/rita.cron
2 * * * * rita import --rolling --logs=/opt/zeek/logs/$(date +\%Y-\%m-\%d)/ --database=goad